# scripts/com.zoeyzhou.timesense.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.zoeyzhou.timesense.backend</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/zoeyzhou/Projects/AI Architect/scripts/start-timesense.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/zoeyzhou/Projects/AI Architect</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/timesense-backend.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/timesense-backend.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/bin:/bin:/usr/sbin:/sbin:/Users/zoeyzhou/Projects/AI Architect/.venv/bin</string>
  </dict>
</dict>
</plist>
